# crontab
SHELL=/bin/bash

0 4 * * * /usr/backup/mariadbBackup.sh > /proc/1/fd/1 2>/proc/1/fd/2
